# systemd/stock-web-v3.service
[Unit]
Description=Stock Web V3 - Real-time IDX Stock Monitoring
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=notify
User=stock_app
Group=stock_app
WorkingDirectory=/opt/stock-web-v3
Environment=PATH=/opt/stock-web-v3/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/stock-web-v3/src
Environment=HOME=/opt/stock-web-v3
EnvironmentFile=/opt/stock-web-v3/.env

ExecStart=/opt/stock-web-v3/venv/bin/uvicorn stock_web_v3.main:app --host 127.0.0.1 --port 8000 --no-access-log
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectHome=false
ProtectSystem=false

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
